# tutorial/toy-test.conf
start on hittoytest
script
    echo "hit-toy-test conf is used" > /tmp/toy-test.log
end script
